# deploy/decks.service
[Unit]
Description=Koho Decks (Presenton fork)

# No Requires=/After=docker.service here: this is a USER unit and
# docker.service is a SYSTEM unit — user systemd can't see it. Docker
# is started by the system before the lingered user session, and if
# it's not running the ExecStart below will exit with a clear error.

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/app
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml up -d --build --remove-orphans production
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml down
TimeoutStartSec=1800
TimeoutStopSec=120

[Install]
WantedBy=default.target
